# staging/systemd/user/task-spooler@.service.d/override.conf
[Service]
Nice=10
CPUSchedulingPolicy=batch
CPUSchedulingResetOnFork=false
IOSchedulingClass=idle
IOSchedulingPriority=7
CPUAffinity=6,7
