# systemd/system/podman/acme-copy-hook.path
[Path]
PathChanged=/var/lib/containers/storage/volumes/acme/_data

[Install]
WantedBy=multi-user.target
