# core/systemd/user/telegram-notification@.service.d/override.conf.sample
[Service]
Environment="TELEGRAM_CHAT_ID=123456789"
Environment="TELEGRAM_BOT_ID=987654321:XXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXX"
