# core/systemd/user/restic-rsync-backup@usb-ssd-polaris-home.service.d/20-override.conf
[Service]
Environment=RESTIC_REPOSITORY=/run/media/zeno/BACKUP/backup/restic/polaris
Environment=RSYNC_TARGET=backup@sirius.fritz.box:/volume1/backup/restic/
Environment=KEYFILE=/home/zeno/.ssh/id_ed25519_rsync
