# systemd/.config/systemd/user/updates.timer
[Unit]
Description=Notify about available updates hourly

[Timer]
OnCalendar=*:16
Persistent=true

[Install]
WantedBy=default.target
