# systemd/.config/systemd/user/updates.service
[Unit]
Description=Notify about available updates

[Service]
ExecStart=/home/zeno/.config/systemd/user/scripts/updates
