# systemd/user/mailsync.service
[Unit]
Description=mailsync

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mailsync
WorkingDirectory=/home/strom/tmpfs