# core/systemd/user/restic-runner-shutdown@.service.d/override.conf
[Service]
Nice=0
CPUSchedulingPolicy=batch
CPUAffinity=
IOSchedulingClass=
IOSchedulingPriority=
