# core/systemd/user/newsboat.service
[Unit]
Description=Fetch newsboat RSS articles

[Service]
Type=oneshot
ExecStart=-/usr/bin/newsboat -x reload

